# crontab
@reboot ~/dotfiles/scripts/misc/change-caps-to-ctrl.sh
