# crontab
@reboot ~/dotfiles/change-caps-to-ctrl.sh
